# src/systemd/ubuntu-transition-tracker.service
[Unit]
Description=Ubuntu Transition Tracker

[Service]
Type=oneshot
# Specifying the user is required, the service errors out without it
User=root
ExecStartPre=/bin/rm -f /srv/transitions/transition-tracker/ben.lock
ExecStartPre=syncmirror
ExecStart=/srv/transitions/transition-tracker/config/go
